# deploy/systemd/ai-lab.service
# Template - deploy/setup.sh generates the actual service file.
# Placeholders: __USER__, __HOME__, __WORKDIR__

[Unit]
Description=ai-lab Personal AI Assistant
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__WORKDIR__/ai-lab
Restart=always
RestartSec=5
Environment=HOME=__HOME__
Environment=PATH=/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
